# etc/rc.local
/usr/local/bin/random-hostname.sh
/usr/bin/openssl rand -hex 16 > /etc/machine-id
