# etc/rc.local
/usr/local/bin/random-hostname.sh
/usr/local/bin/dbus-uuidgen > /etc/machine-id
